# srv.service
[Unit]
Description=Portfolio web server

[Service]
Type=simple
User=exedev
Group=exedev
WorkingDirectory=/home/exedev/portfolio
ExecStart=/home/exedev/portfolio/portfolio
Restart=always
RestartSec=5
Environment=HOME=/home/exedev
Environment=USER=exedev
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
